# freenom-script-renew@.timer
[Unit]
Description="Time freenom-script-renew for %I repetitively"

[Timer]
OnCalender=monthly

[Install]
WantedBy=timers.target
