# systemd/freenom-renew-all.service
[Unit]
Description="Renew all freenom domains"

[Service]
Type=oneshot
ExecStart=/usr/local/bin/freenom.sh -r -a
