# systemd/freenom-renew-all.timer
[Unit]
Description=Weekly renewal of all freenom domains

[Timer]
OnCalendar=weekly

[Install]
WantedBy=timers.target
