# systemd/freenom-update@.timer
[Unit]
Description=Hourly ip update of freenom domain %I

[Timer]
OnBootSec=2min
OnUnitActiveSec=1hour

[Install]
WantedBy=timers.target
