# systemd/freenom-renew@.timer
[Unit]
Description=Weekly renewal of freenom domain %I

[Timer]
OnCalendar=weekly

[Install]
WantedBy=timers.target
