# systemd/freenom-update@.service
[Unit]
Description="Update ip of freenom domain %I"

[Service]
Type=oneshot
ExecStart=/usr/local/bin/freenom.sh -U %I
