# systemd/freenom-renew@.timer
[Unit]
Description="Time freenom-script-renew for %I repetitively"

[Timer]
OnCalendar=weekly

[Install]
WantedBy=timers.target
