# systemd/freenom-script-renew-all@.timer
[Unit]
Description="Time freenom-script-renew all domains repetitively"

[Timer]
OnCalender=weekly

[Install]
WantedBy=timers.target
