# systemd/freenom-renew-all@.timer
[Unit]
Description="Time freenom-script-renew all domains repetitively"

[Timer]
OnCalendar=weekly

[Install]
WantedBy=timers.target
