# systemd/freenom-renew@.service
[Unit]
Description=Renew freenom domain %I

[Service]
Type=oneshot
ExecStart=/usr/local/bin/freenom.sh -r %I
